# roles/vector-map-generator/files/update-vector-maps.service
[Unit]
Description=Update the vector maps

[Service]
Type=simple
ExecStart=/home/chandler/update-vector-maps.sh
